# gunicorn.conf
[program:gunicorn]
command=/usr/bin/gunicorn wsgi -b 127.0.0.1:5000
directory=/deploy/project
